# smart-wifi-manager.service
[Unit]
Description=Smart Wi-Fi Manager
Documentation=https://github.com/alireza787b/smart-wifi-manager
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=__SWM_INSTALL_DIR__/smart-wifi-manager.sh run --config __SWM_CONFIG_PATH__ --status-file __SWM_STATUS_PATH__ --state-dir __SWM_STATE_DIR__ --log-file __SWM_LOG_PATH__
Restart=on-failure
RestartSec=5
RuntimeDirectory=smart-wifi-manager
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=__SWM_CONFIG_DIR__ __SWM_STATE_DIR__ __SWM_RUN_DIR__ __SWM_LOG_DIR__
WorkingDirectory=__SWM_INSTALL_DIR__

[Install]
WantedBy=multi-user.target
